# src/smartscale.socket
[Unit]
Description=FAUST CTF smartscale service

[Socket]
ListenStream=31337
Accept=true

[Install]
WantedBy=sockets.target
